# backend/launchd/com.traveller.backend.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.traveller.backend</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/absoluter/pfad/zu/traveller-charsheet/backend/server.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/absoluter/pfad/zu/traveller-charsheet/backend</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/absoluter/pfad/zu/traveller-charsheet/backend/data/backend.log</string>
    <key>StandardErrorPath</key>
    <string>/absoluter/pfad/zu/traveller-charsheet/backend/data/backend.error.log</string>
</dict>
</plist>
