# backend/launchd/com.traveller.backup.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.traveller.backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/absoluter/pfad/zu/traveller-charsheet/backend/backup-to-git.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/absoluter/pfad/zu/traveller-charsheet/backend</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/absoluter/pfad/zu/traveller-charsheet/backend/data/backup.log</string>
    <key>StandardErrorPath</key>
    <string>/absoluter/pfad/zu/traveller-charsheet/backend/data/backup.error.log</string>
</dict>
</plist>
